# crontab
# m h dom mon dow command
0 0 * * *   /bin/scanner > /tmp/scanner.log 2>&1
